# script/cron.d/root
SHELL=/data/adb/modules/SRauto/bin/bash
0 10 * * 1-5 /data/adb/modules/SRauto/bin/bash "/data/adb/modules/SRauto/script/unlock.sh"
*/5 * * * * /data/adb/modules/SRauto/bin/bash "/data/adb/modules/SRauto/script/log.sh"
